# Cargo/crontab
*/3 * * * * /usr/local/bin/python3 /code/manage.py update_car_locations
